# src/server/worldserver/worldserverd.service.dist
[Unit]
Description=SkyFire World Daemon

[Service]
User=root
WorkingDirectory=/usr/local/skyfire-server/bin
ExecStart=/usr/local/skyfire-server/bin/worldserver

[Install]
WantedBy=multi-user.target

